# install/systemd/symfs@.timer
[Unit]
Description=Periodic run of SymFs for %I

[Timer]
Persistent=true
OnCalendar=daily

[Install]
WantedBy=timers.target
